# deploy/investyo-mcp.service
[Unit]
Description=InvestYo MCP Server (SSE Transport)
Documentation=https://github.com/kevinmarko/Stockpy
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=investyo
Group=investyo
WorkingDirectory=/opt/investyo
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/investyo/.env
ExecStart=/opt/investyo/.venv/bin/python investyo_mcp_server.py --transport sse --port 8080
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/investyo
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=investyo-mcp

[Install]
WantedBy=multi-user.target
