# deploy/investyo-streamlit.service
[Unit]
Description=InvestYo Streamlit Command Center
Documentation=https://github.com/kevinmarko/Stockpy
After=network.target
Wants=investyo-mcp.service

[Service]
Type=simple
User=investyo
Group=investyo
WorkingDirectory=/opt/investyo
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/investyo/.env
ExecStart=/opt/investyo/.venv/bin/python -m streamlit run gui/app.py --server.port 8501 --server.address 127.0.0.1 --server.headless true
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/investyo
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=investyo-streamlit

[Install]
WantedBy=multi-user.target
